# deploy/supervisor/cegu-worker.conf
; =============================================================================
; CEGU pSEO — Supervisor: worker antrian (generate, publish, IndexNow, AI)
; =============================================================================
; Pasang di: /etc/supervisor/conf.d/cegu-worker.conf
; Ganti:
;   - /var/www/cegu   → path root aplikasi Anda
;   - user www-data   → user web server Anda
;
; numprocs=3 → 3 worker paralel menarik dari queue 'default' (job generate &
; publish memecah diri per-chunk, jadi beberapa worker mempercepat produksi
; jutaan halaman). Naikkan sesuai jumlah core CPU.
;
; Setelah menaruh file:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start cegu-worker:*
; =============================================================================

[program:cegu-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cegu/artisan queue:work --sleep=1 --tries=1 --max-time=3600 --max-jobs=1000 --timeout=1800
directory=/var/www/cegu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/cegu/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=1810

; Catatan:
; --max-time=3600 & --max-jobs=1000 → worker restart berkala (cegah memory leak).
; --tries=1 → job berat (generate/publish) tidak diulang membabi buta; kegagalan
;             tercatat agar bisa ditindaklanjuti manual.
; Untuk queue Redis, tidak perlu ubah apa pun di sini (cukup ganti .env).
